# RHEL-9/codeuh/units/neo4j-image.timer
[Unit]
Description=neo4j-image timer to update related podman container image.

[Timer]
OnCalendar=daily
RandomizedDelaySec=900

[Install]
WantedBy=timers.target